# manager/docker/development/crontab
0 * * * * root /usr/local/bin/php /app/bin/console app:weather >> /app/var/log/dev.log
